# Milestone 1/cron/crontab
* * * * * docker compose -f /docker-compose.yml --env-file /.env --project-name ${COMPOSE_PROJECT_NAME} run --rm lego-renew
