# Milestone 1/cron/crontab
* * * * * docker compose -f /docker-compose.yml --env-file /.env --project-name milestone1 run --rm lego-renew
